# roles/website/files/website-extras.timer
[Unit]
Description=Run Website Extras every 2 Minutes

[Timer]
OnBootSec=60s
OnUnitActiveSec=120s

[Install]
WantedBy=timers.target
