# roles/wireguard-concentrator/files/concentrator-route.service
[Unit]
Description=Concentrator Route placeholder

[Service]
Type=simple
ExecStart=true

[Install]
WantedBy=multi-user.target
