# roles/wireguard-concentrator/files/concentrator-route.service
[Unit]
Description=Configure routing for connected nodes

[Service]
Type=simple
ExecStart=/srv/venv/bin/ffbs-concentrator-route
Restart=always
RestartSec=10
StartLimitInterval=0

[Install]
WantedBy=multi-user.target
